# edison/btup.service
[Unit]
Description=BT UP
After=bluetooth.service

[Service]
ExecStart=/bin/bash /opt/diaperpie/btup.sh

[Install]
WantedBy=multi-user.target
